# mage1-fpm/context/etc/profile.d/path.sh
## support packages installed via pip into --user context
export PATH=~/.local/bin:$PATH
